# ops/systemd/invest-certo-pipeline.timer
[Unit]
Description=Run Invest Certo data pipeline every business evening

[Timer]
OnCalendar=Tue..Sat 00:00:00
Persistent=true
Unit=invest-certo-pipeline.service

[Install]
WantedBy=timers.target
